# deploy/systemd/remotetrade-depth-arb.timer
[Unit]
Description=Run RemoteTrade depth arbitrage check every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=10s
Unit=remotetrade-depth-arb.service

[Install]
WantedBy=timers.target
