# deploy/systemd/remotetrade-bitbank-route-probe.timer
[Unit]
Description=Run RemoteTrade bitbank triangular route paper probe every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
AccuracySec=10s
Unit=remotetrade-bitbank-route-probe.service

[Install]
WantedBy=timers.target
