# deploy/systemd/remotetrade-venue-discovery.service
[Unit]
Description=RemoteTrade low-cost venue market discovery
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=remotetrade
Group=remotetrade
WorkingDirectory=/opt/remotetrade
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/remotetrade/.venv/bin/python -m remotetrade.app --discover-venues
