# deploy/systemd/remotetrade-maker-probe.service
[Unit]
Description=RemoteTrade small maker market probe
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=remotetrade
Group=remotetrade
WorkingDirectory=/opt/remotetrade
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/remotetrade/.venv/bin/python -m remotetrade.app --maker-probe
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
