# deploy/systemd/remotetrade-boba-zencha-probe.timer
[Unit]
Description=Run RemoteTrade Boba Zencha StableSwap versus CEX paper probe every 5 minutes

[Timer]
OnCalendar=*-*-* *:01/5:15
AccuracySec=10s
Unit=remotetrade-boba-zencha-probe.service

[Install]
WantedBy=timers.target
