# deploy/systemd/remotetrade-poly-replay.service
[Unit]
Description=RemoteTrade Polymarket BTC 5m replay report
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=remotetrade
Group=remotetrade
WorkingDirectory=/opt/remotetrade
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/remotetrade/.venv/bin/python -m remotetrade.app --poly-replay --discord
