# deploy/systemd/remotetrade-health.timer
[Unit]
Description=Run RemoteTrade health check every 5 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
AccuracySec=30s
Unit=remotetrade-health.service

[Install]
WantedBy=timers.target
