# deploy/systemd/remotetrade-boba-synapse-probe.timer
[Unit]
Description=Run RemoteTrade Boba Synapse canonical stablecoin paper probe every 5 minutes

[Timer]
OnCalendar=*-*-* *:02/5:15
AccuracySec=10s
Unit=remotetrade-boba-synapse-probe.service

[Install]
WantedBy=timers.target
