# deploy/systemd/remotetrade-update.service
[Unit]
Description=Update RemoteTrade from GitHub
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/remotetrade
ExecStart=/usr/bin/bash /opt/remotetrade/deploy/update.sh
