# deploy/systemd/remotetrade-poly-stock-bridge.timer
[Unit]
Description=Collect Polymarket-to-stock bridge observations every 10 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=10min
AccuracySec=15s
Unit=remotetrade-poly-stock-bridge.service

[Install]
WantedBy=timers.target
